# templates/ttySX.conf.erb
# ttyS1 - getty
#
# This service maintains a getty on ttyS1 from the point the system is
# started until it is shut down again, providing Serial-Over-Lan access
# from the BMC of the xenon supermicro machines.

start on runlevel [23] and not-container

stop on runlevel [!23]

respawn
exec /sbin/getty -L <%= serial_tty %> 115200 vt102
